# roles/restake/files/restake.timer
[Unit]
Description=Restake bot timer

[Timer]
AccuracySec=1min
OnCalendar=*-*-* 01:00:00

[Install]
WantedBy=timers.target
